# apache.sh
#!/bin/bash
### BEGIN INIT INFO
# Provides:          apache.sh
# Required-Start:    apache2 start
# Required-Stop:     apache2 start
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Simple configuration procedure
# Description:       A simple script to start the needed configuration procedure
### END INIT INFO

# Creates a symbolic link to the 

